# software/boost.service
[Unit]
Description=Boost
 
[Service]
ExecStart=/home/giatrung2012/software/gboost.sh
 
[Install]
WantedBy=multi-user.target
